# rootfs/etc/services.d/mlat-client/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

APPNAME="$(hostname)/mlat-client"
echo "[$APPNAME][$(date)] Started as an s6 service"

MLAT_CMD="/usr/bin/mlat-client"
MLAT_FLAG="/run/mlat-down"
RESTARTTIMER=10

[[ "${MLAT_INPUT_TYPE}" == "" ]] && MLAT_INPUT_TYPE="auto" || true
[[ "${MLAT_RESULTS}" == "" ]] && MLAT_RESULTS="beast,listen,30105" || true


# Sleep forever if the initialization process found that MLAT cannot be started
[[ -f $MLAT_FLAG ]] && sleep infinity || true

MLAT_PARAM=(--input-type "${MLAT_INPUT_TYPE}")
# MLAT_PARAM+=(--no-udp)
MLAT_PARAM+=(--input-connect "${MLAT_HOST}")
MLAT_PARAM+=(--server "${MLAT_SERVER}")
MLAT_PARAM+=(--lat "${LAT}")
MLAT_PARAM+=(--lon "${LON}")
MLAT_PARAM+=(--alt "${ALT}")
MLAT_PARAM+=(--results "${MLAT_RESULTS}")
MLAT_PARAM+=(--user "${FEEDER_KEY}")

[[ -n "${PRIVATE_MLAT}" ]] && MLAT_PARAM+=(--privacy) || true

set -eo pipefail
[[ "$VERBOSE" == "ON" ]] && echo "[$APPNAME][$(date)] Starting MLAT with: ${MLAT_CMD} ${MLAT_PARAM[@]}" || true
[[ "$VERBOSE" != "ON" ]] && exec "${MLAT_CMD}" "${MLAT_PARAM[@]}" 2>&1 1>/dev/null | stdbuf -o0 awk -v app="$APPNAME" '{print "[" app "] " $0}' || true
if [[ "$VERBOSE" == "ON" ]]
then
   exec "${MLAT_CMD}" "${MLAT_PARAM[@]}" 2>&1 |  stdbuf -o0 awk -v app="$APPNAME" -v date="$(date)" '{print "[" app "][" date "] " $0}'
   a=$?
else
  exec "${MLAT_CMD}" "${MLAT_PARAM[@]}" >/dev/null 2>&1 |  stdbuf -o0 awk -v app="$APPNAME" -v date="$(date)" '{print "[" app "][" date "] " $0}'
  a=$?
fi

echo "[$APPNAME][$(date)] MLAT-client exited with exit code $a. Attempting to restart in $RESTARTTIMER seconds"
sleep $RESTARTTIMER
