# rootfs/etc/services.d/mlat-client/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

MLAT_CMD="/usr/binmlat-client"

[[ "${MLAT_INPUT_TYPE}" == "" ]] && MLAT_INPUT_TYPE=dump1090

MLAT_PARAM=(--input-type "${MLAT_INPUT_TYPE}")
# MLAT_PARAM+=(--no-udp)
MLAT_PARAM+=(--input-connect "${MLAT_HOST}")
MLAT_PARAM+=(--server "${MLAT_SERVER}")
MLAT_PARAM+=(--lat "${LAT}")
MLAT_PARAM+=(--lon "${LON}")
MLAT_PARAM+=(--alt "${ALT}")
MLAT_PARAM+=(--results beast,listen,50010)
MLAT_PARAM+=(--user "${FEEDER_KEY}")
if [[ ${PRIVATE_MLAT} == true ]]; then
    MLAT_PARAM+=(--privacy)
fi

set -eo pipefail
echo "Starting MLAT with: ${MLAT_CMD} ${MLAT_PARAM[@]}"
exec \
    "${MLAT_CMD}" "${MLAT_PARAM[@]}" 2>&1 | awk -W interactive '{print "[mlat-client] " $0}'
