# rootfs/etc/services.d/radarvirtuel/run
#!/usr/bin/with-contenv bash
#shellcheck shell=bash

APPNAME="$(hostname)/radarvirtuel"
[[ "$RV_SERVER" == "" ]] && RV_SERVER="mg2.adsbnetwork.com:50046"
[[ "$VERBOSE" == "ON" ]] && verbose="-v" || verbose=""

echo "[$APPNAME][$(date)] RadarVirtuel was started as an s6 service"
echo "[$APPNAME][$(date)] Feeder key: $FEEDER_KEY"
echo "[$APPNAME][$(date)] Source host: $SOURCE_HOST"
echo "[$APPNAME][$(date)] RadarVirtuel Server: $RV_SERVER"
echo "[$APPNAME][$(date)] Verbose: $([[ "$VERBOSE" == "ON" ]] && echo "ON" || echo "OFF")"
/home/py/ANfeeder $verbose -i $FEEDER_KEY -d $RV_SERVER -s $SOURCE_HOST
a=$?
echo "[$APPNAME][$(date)] RadarVirtuel exited with exit code $a"
